# name.timer
[Timer]
OnStartupSec=5m
# retry every hour since service checks whether to actually run
OnActiveSec=1h
